# kubernetes/apiserver/kube-apiserver.service
[Unit]
Description=Kubernetes API Server
Documentation=https://github.com/kubernetes/kubernetes

[Service]
Type=simple
ExecStart=/bin/bash -c "GOMAXPROCS=$(nproc) /usr/local/programes/kubernetes/services/apiserver/main.sh"

[Install]
WantedBy=multi-user.target