# kubernetes/controller-manager/kube-controller-manager.service
[Unit]
Description=Kubernetes controller manager
Documentation=https://github.com/kubernetes/kubernetes

[Service]
Type=simple
ExecStart=/bin/bash -c "GOMAXPROCS=$(nproc) /usr/local/programes/kubernetes/services/controller-manager/main.sh"

[Install]
WantedBy=multi-user.target